# utils/supervisor/salt-celery-test.ini
[program:celery-test]
command=celery -A salts_prj worker -l info -c 3 --workdir /data/qa/ltbot/salts --logfile /tmp/celery-salts-test.log
directory=/data/qa/ltbot/venv
user=uwsgi
autostart=true
autorestart=true
redirect_stderr=true
