# utils/supervisor/salt-pypi-server.ini
[program:pypi-server]
command=/usr/bin/python /usr/bin/pypi-server -p 8080 -P /data/qa/pypiserver/.htaccess -o --overwrite /data/qa/pypiserver/packages
directory=/data/qa/ltbot/venv
user="<username>"
autostart=true
autorestart=true
redirect_stderr=true
