# utils/supervisor/salt-celery-prod.ini
[program:celery-prod]
command=celery -A salts_prj worker -l info -c 5 --workdir /usr/local/salts --logfile /tmp/celery-salts-prod.log
directory=/data/qa/ltbot/venv
user=uwsgi
autostart=true
autorestart=true
redirect_stderr=true
