# utils/supervisor/tank-salts-tank.ini
[program:salts-tank]
command=/data/qa/ltbot/venv/bin/run-salts-tank
directory=/data/qa/ltbot/venv
user="<username>"
autostart=true
autorestart=true
redirect_stderr=true
