# utils/deploy-yatank-env/supervisord.d/yandex-tank-api-server.ini
[program:yatankapi]
command=/data/qa/ltbot/venv/bin/yandextankapiserver.sh
directory=/data/qa/ltbot/venv
user="<username>"
autostart=true
autorestart=true
redirect_stderr=true
